# ai-web-app.service
[Unit]
Description=AI Web Application Stack
After=docker.service
Requires=docker.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/ai-app  
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
TimeoutStartSec=300
TimeoutStopSec=30

User=tesfaye
Group=tesfaye  # ← CHANGED from docker to tesfaye

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/ai-app  # ← CHANGED from ai-web-app to ai-app

[Install]
WantedBy=multi-user.target
